# .devcontainer/conf.d/ursim.conf
[program:ursim]
command=/ursim/start-ursim.sh
redirect_stderr=true
autorestart=true
