# .devcontainer/conf.d/xvfb.conf
[program:xvfb]
#command=Xvfb :1 -screen 0 "%(ENV_DISPLAY_WIDTH)s"x"%(ENV_DISPLAY_HEIGHT)s"x16 -listen tcp -ac
command=Xvfb :1 -screen 0 1280x800x16
autorestart=true
